# checkrmd.service
[Service]
ExecPreStart="/usr/bin/sleep 15"
ExecStart="/usr/bin/checkrmd"
Restart=always

